# systemd/mowerbase-gps.service
[Unit]
Description=BaseStation GPS Service — LC29H(BS) serial reader and Survey-In state machine
Documentation=https://github.com/yourusername/mowerbase
After=network.target
Wants=network.target

[Service]
Type=simple
User=mowerbase
Group=mowerbase

# systemd creates and owns these directories with correct permissions before start
RuntimeDirectory=mowerbase
RuntimeDirectoryMode=0755
LogsDirectory=mowerbase
LogsDirectoryMode=0755

ExecStart=/opt/mowerbase/venv/bin/python /opt/mowerbase/gps.py

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Allow access to serial port and GPIO
SupplementaryGroups=dialout gpio

# PID for signal targeting
PIDFile=/run/mowerbase/gps.pid

[Install]
WantedBy=multi-user.target
